# sv/tty5/run
#!/bin/execlineb -P
if { pipeline { s6-rc -ba list } grep -qFx mount-filesystems }
envfile /etc/s6-rc/conf.d/tty5.conf
importas -uD "yes" SPAWN SPAWN
importas -sCuD "" ARGS ARGS
importas -sCuD "agetty" GETTY GETTY
if -t { test -e /dev/tty5 }
if -t { test ${SPAWN} = "yes" }
	exec ${GETTY} -L -8 tty5 115200 ${ARGS}
