# sv/tmpfiles-setup/up
#!/bin/execlineb -P
exec systemd-tmpfiles --exclude-prefix=/dev --create --remove --boot
