# sv/locale/up
#!/bin/execlineb -P
exec sh /etc/profile.d/locale.sh
