# sv/nix-daemon-log/run
#!/bin/execlineb -P
s6-envdir /etc/s6-rc/conf.d/nix-daemon.conf.d
importas -sCuD "n3 s2000000 T" DIRECTIVES DIRECTIVES
ifelse { test -w /var/log } {
	foreground { install -d -o s6log -g s6log /var/log/nix-daemon }
	s6-setuidgid s6log exec -c s6-log -d3 -b -- ${DIRECTIVES} /var/log/nix-daemon
}
foreground { install -d -o s6log -g s6log /run/log/nix-daemon }
s6-setuidgid s6log exec -c s6-log -d3 -b -- ${DIRECTIVES} /run/log/nix-daemon
