# sv/mount-filesystems/up
#!/bin/execlineb -P
foreground { mount -a -O no_netdev }
