# sv/binfmt/up
#!/bin/execlineb -P
exec sh /etc/s6/sv/binfmt/shell_up
