# sv/wpa_supplicant-srv/run
#!/bin/execlineb -P
s6-envdir /etc/s6-rc/conf.d/wpa_supplicant.conf.d
importas -uD "wlan0" WPA_INTERFACE WPA_INTERFACE
fdmove -c 2 1
exec /usr/bin/wpa_supplicant -i ${WPA_INTERFACE} -c /etc/wpa_supplicant/wpa_supplicant.conf.d
