# sv/modules/up
#!/bin/execlineb -P
foreground { modules-load }
