# sv/syslog-ng-log/run
#!/bin/execlineb -P
envfile /etc/s6-rc/conf.d/syslog-ng.conf
importas -sCuD "n3 s2000000 T" DIRECTIVES DIRECTIVES
ifelse { test -w /var/log } {
	foreground { install -d -o s6log -g s6log /var/log/syslog-ng }
	s6-setuidgid s6log exec -c s6-log -d3 -b -- ${DIRECTIVES} /var/log/syslog-ng
}
foreground { install -d -o s6log -g s6log /run/log/syslog-ng }
s6-setuidgid s6log exec -c s6-log -d3 -b -- ${DIRECTIVES} /run/log/syslog-ng
