# sv/pcscd-log/run
#!/bin/execlineb -P
s6-envdir /etc/s6-rc/conf.d/pcsclite.conf.d
importas -sCuD "n3 s2000000 T" DIRECTIVES DIRECTIVES
ifelse { test -w /var/log } {
	foreground { install -d -o s6log -g s6log /var/log/pcscd }
	s6-setuidgid s6log exec -c s6-log -d3 -b -- ${DIRECTIVES} /var/log/pcscd
}
foreground { install -d -o s6log -g s6log /run/log/pcscd }
s6-setuidgid s6log exec -c s6-log -d3 -b -- ${DIRECTIVES} /run/log/pcscd
