# sv/thinkfan-srv/run
#!/usr/bin/execlineb -P
fdmove -c 2 1
/usr/bin/thinkfan -q -n -s5 -c /etc/thinkfan.conf
