# sv/remount-root/up
#!/bin/execlineb -P
foreground { mount -o remount,rw / }
foreground { sh /etc/s6/log-service-reload.sh }
