# sv/hwclock/up
#!/bin/execlineb -P
s6-envdir /etc/s6-rc/conf.d/hwclock.conf.d
importas -uD "UTC" HARDWARECLOCK HARDWARECLOCK
ifelse { test $HARDWARECLOCK = UTC } { hwclock --systz --utc --noadjfile } hwclock --systz --localtime --noadjfile
