# sv/hwclock/up
#!/bin/execlineb -P
envfile /etc/s6-rc/conf.d/hwclock.conf
importas -uD "UTC" HARDWARECLOCK HARDWARECLOCK
ifelse { test $HARDWARECLOCK = UTC } { hwclock --systz --utc --noadjfile } hwclock --systz --localtime --noadjfile
