# sv/tmpfiles-dev/up
#!/bin/execlineb -P
exec /usr/bin/systemd-tmpfiles --prefix=/dev --create --boot
