# sv/network-detection/up
#!/bin/execlineb -P
background { sh /etc/s6/sv/network-detection/shell_up }
