# sv/nix-daemon-srv/run
#!/bin/execlineb -P
fdmove -c 2 1 /nix/var/nix/profiles/default/bin/nix-daemon 
