# sv/mount-cgroups/up
#!/bin/execlineb -P
s6-envdir /etc/s6-rc/conf.d/mount-cgroups.conf.d
importas -uD "hybrid" CGROUP_MODE CGROUP_MODE
importas -uD "none" CGROUP_CONTROLLERS CGROUP_CONTROLLERS
importas -uD "true" HAVE_CONTROLLER1_GROUPS HAVE_CONTROLLER1_GROUPS
exec sh /etc/s6/sv/mount-cgroups/shell_up $CGROUP_MODE $CGROUP_CONTROLLERS $HAVE_CONTROLLER1_GROUPS
